# ansible/roboshop/roles/cart/templates/systemd.service
[Unit]
Description = Cart Service
[Service]
User=roboshop
Environment=REDIS_HOST=redis-{{ENV}}.knowaws.com
Environment=CATALOGUE_HOST=catalogue-{{ENV}}.knowaws.com
ExecStart=/bin/node /home/roboshop/cart/server.js
SyslogIdentifier=cart

[Install]
WantedBy=multi-user.target